# frontend/deploy/conf/supervisor/firstrun.conf
[program:firstrun]
command=/bin/bash /tmp/firstrun.sh
user=root
autorestart=false
startretries=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/firstrun.log
chmod=0755