# frontend/deploy/conf/supervisor/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
