# api/deploy/conf/supervisor/cron.conf
[program:cron]
command=cron -f -L 15
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cron.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
